# msm-agent/msm-agent.service.template
[Unit]
Description=MSM Agent (Node Runtime)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=msm
Group=msm
WorkingDirectory={{AGENT_DIR}}
# Systemd units do not inherit login PATH. venv first, then system paths.
Environment="PATH={{AGENT_DIR}}/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="DOCKER_HOST={{DOCKER_HOST}}"
Environment="XDG_RUNTIME_DIR=/run/user/{{MSM_UID}}"
EnvironmentFile=-{{AGENT_DIR}}/.env
ExecStart={{AGENT_DIR}}/venv/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Agent only needs servers data + rootless docker runtime socket + TLS certs
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
StateDirectory=msm-agent
ReadWritePaths={{AGENT_DIR}} {{SERVERS_DIR}} /run/user/{{MSM_UID}} /var/lib/msm-agent

[Install]
WantedBy=multi-user.target
